# systemd/user/dotfiles-maintenance.service
[Unit]
Description=Dotfiles Maintenance Check

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -lc '%h/.local/bin/dotfiles-maintenance'
